# ansible/roles/openclaw/files/clawless-wake-greet.service
[Unit]
Description=Clawless wake greeting (proactive message on provision or resume)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clawless-wake-greet
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
